# etc/service/30redis/run
#!/bin/bash
cd /usr/bin/ && exec /usr/bin/redis-server /etc/redis/redis.conf
