# etc/service/80websocket/run
#!/bin/bash
[ -r /etc/default/diaspora ] && . /etc/default/diaspora
wait_service redis
cd $DIASPORA_HOME && 
exec /usr/bin/ruby $DIASPORA_HOME/script/websocket_server.rb
