# etc/service/10mongo/run
#!/bin/bash
exec /usr/bin/mongod -f /etc/mongodb.conf
