# etc/service/90resque_worker/run
#!/bin/bash
[ -r /etc/default/diaspora ] && . /etc/default/diaspora
wait_service websocket
cd $DIASPORA_HOME && QUEUE=* HOME=/root exec /usr/local/bin/bundle exec rake resque:work
