# etc/service/99thin_1/run
#!/bin/bash
[ -r /etc/default/diaspora ] && . /etc/default/diaspora
# Expects to be in a directory that ends with _[0-9], ie thin_0...
exec_thin 
